# new_code/deploy/systemd/veyrasoul-v2-cloudflared.service
[Unit]
Description=VeyraSoul V2 Cloudflare Tunnel
After=network-online.target veyrasoul-v2.service
Wants=network-online.target
Requires=veyrasoul-v2.service
ConditionPathExists=/etc/cloudflared/token

[Service]
Type=simple
ExecStart=/usr/bin/cloudflared tunnel --no-autoupdate --protocol auto --edge-ip-version 4 --url http://127.0.0.1:8875 run --token-file /etc/cloudflared/token
Restart=always
RestartSec=5
TimeoutStartSec=90
TimeoutStopSec=20
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6
LockPersonality=true
MemoryDenyWriteExecute=true
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
